# infrastructure/modules/ec2-monitoring/scripts/ssh-key-checker.timer
[Unit]
Description=Run SSH Key Format Checker periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
